# systemd/ffxiv-bot.service
[Unit]
Description=FFXIV Character Management Discord Bot
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=ffxiv-bot
Group=ffxiv-bot
WorkingDirectory=/opt/ffxiv-discord-bot
Environment="PYTHONPATH=/opt/ffxiv-discord-bot"
Environment="PATH=/opt/ffxiv-discord-bot/venv/bin:$PATH"
Environment="VIRTUAL_ENV=/opt/ffxiv-discord-bot/venv"

# Use the Python from the virtual environment
ExecStart=/opt/ffxiv-discord-bot/venv/bin/python /opt/ffxiv-discord-bot/bot.py

# Restart settings
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Memory management (update MemoryLimit to MemoryMax as suggested in logs)
MemoryMax=512M

[Install]
WantedBy=multi-user.target